# S7Logger.service
[Unit]
Description=S7Logger service app

[Service]
Type=notify
ExecStart=/opt/S7Logger/S7Logger

[Install]
WantedBy=multi-user.target
